# terraform/vm_provisioning/configs/traefik.service
# vm_provisioning/configs/traefik.service
[Unit]
Description=Traefik Edge Router
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/traefik --configfile=/etc/traefik/traefik.yml
Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target